# vendor/systemd/db-migrate.service
[Unit]
Description=Matreon migrate database
After=docker-compose-matreon.service

[Service]
Restart=no

WorkingDirectory=/usr/local/src/matreon
EnvironmentFile=/matreon_data/.env

ExecStart=/usr/bin/docker-compose run web rake db:migrate

[Install]
WantedBy=multi-user.target
