# vendor/systemd/rails.service
[Unit]
Description=Prepares and starts the Matreon Rails server
After=postgresql-10.service

[Service]
WorkingDirectory=/var/www/matreon
EnvironmentFile=/home/matreon/.env

User=matreon
Type=simple
Restart=on-failure
PIDFile=/var/www/matreon/tmp/puma.pid

ExecStartPre=/usr/local/bin/bundle --without development:test
ExecStartPre=/usr/local/bin/bundle exec rake db:migrate
ExecStartPre=/usr/bin/yarn
ExecStartPre=/usr/local/bin/bundle exec rake assets:precompile

ExecStart=/usr/local/bin/bundle exec /home/matreon/bin/puma -C config/puma.rb

[Install]
WantedBy=multi-user.target
