# vendor/bitcoin/ibd-shutdown.path
[Path]
PathExists=/home/bitcoin/.ibd_service_requests_shutdown
