# vendor/bitcoin/initial-blockchain-download.service
[Unit]
Description=Check if initial blockchain download has done and if we're on a special machine.
After=network.target
ConditionPathExists=!/home/bitcoin/.ibd_service_finished

# Newer versions of systemd (>~ october 2017)
# StartLimitIntervalSec=10
# StartLimitBurst=1
# RestartSec=60

[Service]
ExecStart=/usr/local/src/matreon/vendor/bitcoin/initial_blockchain_download_and_prune.sh  

WorkingDirectory=/home/bitcoin
User=bitcoin
Type=simple
Restart=on-failure

StartLimitInterval=10
StartLimitBurst=1
RestartSec=60

# Hardening measures
####################

# Provide a private /tmp and /var/tmp.
PrivateTmp=true

# Mount /usr, /boot/ and /etc read-only for the process.
ProtectSystem=full

# Disallow the process and all of its children to gain
# new privileges through execve().
NoNewPrivileges=true

# Use a new /dev namespace only populated with API pseudo devices
# such as /dev/null, /dev/zero and /dev/random.
PrivateDevices=true

[Install]
WantedBy=multi-user.target
