# vendor/charge/lightning-charge.path
[Path]
PathExists=/home/bitcoin/.ibd_service_finished
